# files/system/usr/lib/systemd/system/brew-update.service.d/upgrade.conf
[Service]
ExecStart=
ExecStart=/usr/bin/bash -c "/home/linuxbrew/.linuxbrew/bin/brew update && /home/linuxbrew/.linuxbrew/bin/brew upgrade"
