# files/system/usr/share/samuel-niri/dotfiles/dot_config/systemd/user/battery-notify.service
[Unit]
Description=Battery Low Notification Service

[Service]
ExecStart=%h/.local/bin/battery-notify.sh
Restart=on-failure
RestartSec=60

[Install]
WantedBy=default.target
