# files/system/usr/share/samuel-niri/dotfiles/dot_config/systemd/user/cliphist-wipe.service
[Unit]
Description=Wipe cliphist clipboard history on boot

[Service]
Type=oneshot
ExecStart=/usr/bin/cliphist wipe

[Install]
WantedBy=default.target
