# files/system/usr/share/samuel-niri/dotfiles/dot_config/systemd/user/empty-trash.service
[Unit]
Description=Empty Trash on Login
After=default.target

[Service]
Type=oneshot
ExecStart=/usr/bin/bash -c "yes | /usr/bin/gio trash --empty"

[Install]
WantedBy=default.target